# bell-ringer.service
[Unit]
Description=Ring the bell

[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart=sudo /usr/local/bell-ringer/bellserver.py &> /tmp/bell_log.log 

[Install]
WantedBy=multi-user.target

